# deploy/insight-api.service
[Unit]
Description=INSIGHT.AI FastAPI backend
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/insight-ai/app
EnvironmentFile=/opt/insight-ai/.env
ExecStart=/opt/insight-ai/app/venv/bin/uvicorn functions.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
